# modules/raspi/files/epaperdisplay/epim.service
[Unit]
Description=EPIM: EPaper Image Maker

[Service]
Type=oneshot
ExecStart=/home/woolie/epim/draw
User=woolie
Group=woolie
TimeoutStartSec=60s
